# scripts/veilguard-tcmm-stats.service
[Unit]
Description=TCMM + LanceDB performance snapshot
# Doesn't depend on TCMM being up — the script handles "service down"
# gracefully and logs healthy=false when it can't connect.

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/rudol/veilguard/scripts/tcmm_stats.py
# Quick run; bail if it ever hangs (network timeout to /health is 4s,
# du can be slow on cold cache for 2 GB tables — give 60s headroom).
TimeoutStartSec=60
# stdout already goes to the JSONL log via the script. Keep stderr in
# journal for debugging if the script itself crashes.
StandardOutput=null
StandardError=journal
