# scripts/veilguard-tcmm-maintenance.service
[Unit]
Description=LanceDB compaction + version cleanup for TCMM
# Run AFTER TCMM is set up; we'll stop it ourselves in ExecStartPre.
After=veilguard-tcmm.service

[Service]
Type=oneshot
# Run as root so systemctl stop/start works without password prompt,
# but the actual Python script runs AS rudol because Lance is installed
# in rudol's pip user-site (~/.local/lib/python*/site-packages). If we
# ran the Python as root we'd get ModuleNotFoundError: No module named
# 'lance' — bit us once during the first install.
ExecStartPre=/bin/systemctl stop veilguard-tcmm
ExecStart=/usr/bin/sudo -u rudol /usr/bin/python3 /home/rudol/veilguard/scripts/lance_maintenance.py
# Always restart TCMM, even if maintenance failed.
ExecStartPost=/bin/systemctl start veilguard-tcmm
StandardOutput=append:/var/log/veilguard-lance-maintenance.log
StandardError=append:/var/log/veilguard-lance-maintenance.log
# Don't run forever if Lance hangs.
TimeoutStartSec=900
